# cron/crontabs/root
30 * * * *    /usr/bin/ntpdate -s -u pool.ntp.org
